# device-mail-el9/src/main/resources/usr/lib/systemd/system/mail-autodiscovery.target
# systemd service for managing all mail autoconfiguration.

[Unit]
Description=Mail Configuration Autodiscovery

[Install]
WantedBy=multi-user.target

